# ops/autopilot/voicetranslator-autopilot-stop-all.service
[Unit]
Description=Stop incomplete VoiceTranslator Jules sessions
Wants=network-online.target
After=network-online.target
ConditionPathExists=/etc/voicetranslator-autopilot.env

[Service]
Type=oneshot
User=root
Group=root
EnvironmentFile=/etc/voicetranslator-autopilot.env
ExecStart=/usr/bin/python3 /opt/voicetranslator-autopilot/autopilot.py stop-all-sessions
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=full
ReadWritePaths=/opt/voicetranslator-autopilot /var/lib/voicetranslator-autopilot
